# config/com.gmail-assistant.process-inbox.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gmail-assistant.process-inbox</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/tomas/git/ai/gmail-assistant/bin/process-inbox</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/tomas/git/ai/gmail-assistant</string>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/tomas/.local/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/tomas/git/ai/gmail-assistant/logs/process-inbox.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tomas/git/ai/gmail-assistant/logs/process-inbox-error.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
